# daily_audit.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=acline@us.ibm.com
HOME=/root
0 12 * * *  root   cd /opt/vadar; ruby vadar.rb daily_audit | sendmail -t
