# quarterly_report.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=acline@us.ibm.com
HOME=/root
0 12 1 1,4,7,10 *  root   cd /opt/vadar; ruby vadar.rb quarterly_report | sendmail -t
